# deploy/reports/spendhound-monthly-reports.timer
[Unit]
Description=Run SpendHound monthly reports on the first day of each month

[Timer]
OnCalendar=*-*-01 00:20:00
RandomizedDelaySec=20m
Persistent=true
Unit=spendhound-monthly-reports.service

[Install]
WantedBy=timers.target
